# systemd/dynipset.target
[Unit]
Description=DynIPSet target unit
StopWhenUnneeded=true
